# centos6_Yum_Apache_php5.6/supervisord.conf
[supervisord]
nodaemon=true

#[program:sshd]
#command=/usr/sbin/sshd -D
#[program:mysqld]
#command=/usr/bin/mysqld_safe

#[program:nginx] 
#command = /usr/sbin/nginx  -g 'daemon off;'

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND